# roles/mailserver/files/export/systemd.service
[Unit]
Description=Exports certificates for mailserver
After=network.target

[Service]
Type=oneshot
ExecStart=/var/srv/mailserver/export.sh
