# roles/mailserver/files/export/systemd.timer
[Unit]
Description=Run export.sh once a day

[Timer]
OnCalendar=daily
Persistent=true
Unit=export.service
